# utils/launchd/de.gsi.odc.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
	<dict>
		<key>Label</key>
		<string>de.gsi.dds-control-server</string>
		<key>Program</key>
		<string>@CMAKE_INSTALL_PREFIX@/bin/dds-control-server</string>
		<key>StandardOutPath</key>
		<string>/tmp/dds-control-server_out.log</string>
		<key>StandardErrorPath</key>
		<string>/tmp/dds-control-server_err.log</string>
	</dict>
</plist>
