# utils/systemd/odc.service.in
[Unit]
Description=ODC server

[Service]
Type=simple
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/odc-grpc-server

[Install]
WantedBy=multi-user.target
