# utils/systemd/odc.service.in
[Unit]
Description=DDS control server

[Service]
Type=simple
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/dds-control-server

[Install]
WantedBy=multi-user.target
